# service/my_go_cron.service
[Unit]
Description=My Go Program

[Service]
ExecStart=main
Restart=always
User=rishavmngo
Group=rishavmngo

[Install]
WantedBy=mutli-user.target
